# nc-lightsource-opcua.service
[Unit]
Description=NectarCAM Calibration Light Source OPC UA Server
Documentation=https://gitlab.cta-observatory.org/cta-array-elements/mst/nectarcam
# Wait for the network stack to be fully up before starting, so the
# TCP connection to the device and the OPC UA endpoint binding both succeed.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# --- Adjust these to match your deployment ---
User=nectarcam
Group=nectarcam
WorkingDirectory=/opt/nectarcam

ExecStart=/opt/nectarcam/venv/bin/python \
    /opt/nectarcam/nc_lightsource_asyncua_server.py \
    --product FF \
    --address 10.11.4.69 \
    --opcua-endpoint opc.tcp://0.0.0.0:4840/nectarcam/ \
    --opcua-namespace http://cta-observatory.org/nectarcam/1/calibrationlightsource/FF \
    --auto-reconnect \
    --poll-interval 1.0 \
    --log-level INFO

# systemd sends SIGTERM on stop; the server catches it and shuts down cleanly.
# Allow up to 10 s for the OPC UA session teardown before SIGKILL.
TimeoutStopSec=10

# Restart automatically on crash (exit code != 0); not on clean exit (code 0)
# or deliberate stop via systemctl stop.
Restart=on-failure
RestartSec=10

# Journal integration – timestamps and log level are added by journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nc-lightsource-opcua

[Install]
WantedBy=multi-user.target
